# dist/init/systemd.service
[Unit]
Description=Ishtar
After=network.target

[Service]
User=www-data
Restart=always
Type=simple
WorkingDirectory=/opt/ishtar
ExecStart=/opt/ishtar/bin/ishtar-linux64

[Install]
WantedBy=default.target

